# farm/crontab
@reboot root    bash /opt/chia-scripts/farm/clean.sh
@reboot klack   bash /opt/chia-scripts/farm/farm.sh
@reboot klack   bash /opt/chia-scripts/farm/start_hud.sh
