# farm/crontab
@reboot root    bash /opt/chia-scripts/farm/clean.sh
@reboot root    su -c 'sudo /opt/chia-scripts/farm/farm.sh' user
@reboot user    su -c 'sudo /opt/chia-scripts/farm/start_hud.sh' user
