# farm/crontab
@reboot root    bash /opt/chia-scripts/farm/clean.sh
@reboot root    su -c 'sudo /opt/chia-scripts/farm/farm.sh' user
@reboot root    screen -dm -c /opt/chia-scripts/farm/hud.screen
