# plot/crontab
@reboot root    bash /opt/chia-scripts/plot/clean.sh
@reboot root    bash /opt/chia-scripts/plot/start_pm.sh