# plot/crontab
@reboot root    bash /opt/chia-scripts/plot/clean.sh
@reboot root    screen -dm -c /opt/chia-scripts/plot/start.screen
