# plotter_daemon/crontab
@reboot root    bash /opt/chia-scripts/plotter_daemon/clean.sh
@reboot root    screen -dm -c /opt/chia-scripts/plotter_daemon/plotter_hud.screen